# jenkins-docker-builder/supervisord.conf
[supervisord]
nodaemon=true

[include]
files = /etc/supervisor/conf.d/*.conf

[program:docker]
command=/wrapdocker

[program:jenkins-slave]
command=exec /start-java-builder.sh {{swarm_params}}
